# src/cli/startup/openrc-service
#!/sbin/openrc-run

description="{description}"
command="{exe_path}"
command_args="{args}"
pidfile="/run/${{RC_SVCNAME}}.pid"
name="hptp"

depend() {{
    need net
}}

start() {{
    ebegin "Starting ${{name}}"
    start-stop-daemon --start --pidfile ${{pidfile}} --background --exec ${{command}} -- ${{command_args}}
    eend $?
}}

stop() {{
    ebegin "Stopping ${{name}}"
    start-stop-daemon --stop --pidfile ${{pidfile}}
    eend $?
}}